# wgmon.service
[Unit]
Description=WireGuard auto-connect daemon based on WiFi SSID (via iwd)
After=iwd.service network.target
Wants=iwd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/wgmon
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wgmon

# Hardening — wg-quick requires root for network interface management,
# but we can restrict what the process can access beyond that.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
LockPersonality=yes
RestrictSUIDSGID=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictNamespaces=~user mnt pid ipc uts
# Allow read access to WireGuard and wgmon configs.
ReadOnlyPaths=/etc/wireguard /etc/wgmon

[Install]
WantedBy=multi-user.target
